# kanata/.config/autostart/Kanata.desktop
[Desktop Entry]
Exec=/home/tobi/.nix-profile/bin/kanata
Name=Kanata
Type=Application
Version=1.0
