# kanata/.config/autostart/kanata.desktop
[Desktop Entry]
Exec=/bin/bash -c "$HOME/.nix-profile/bin/kanata --nodelay"
Name=Kanata
Type=Application
Version=1.0
