# kanata/.config/autostart/kanata.desktop
[Desktop Entry]
Exec=$HOME/.nix-profile/bin/kanata --nodelay
Name=Kanata
Type=Application
Version=1.0
