# kanata/.config/autostart/Kanata.desktop
[Desktop Entry]
Exec=~/.nix-profile/bin/kanata
Name=Kanata
Type=Application
Version=1.0
